# deploy/kuwai-agent.timer
# /etc/systemd/system/kuwai-agent.timer
# UTC OnCalendar, Persistent=true (a fire missed during VPS downtime runs once
# on recovery; cron silently drops missed fires, disqualifying). The runtime
# clock.is_open gate is the second layer; no equities order ever fires outside
# the sealed RTH guard. ~30-min cadence across the US regular session plus a
# pre-close wake (the run self-guards the exact window from the broker
# calendar API; the timer is intentionally a superset in UTC).
[Unit]
Description=KUWAI commitment-v2 wake schedule (UTC, Persistent)

[Timer]
# 13:30-20:00 UTC covers 09:30-16:00 ET (DST-shift tolerant; the run derives
# the precise open/close/half-day from the Alpaca calendar API at runtime).
# One OnCalendar= per fire (systemd accumulates repeated lines; a combined
# comma time-list is not valid calendar syntax).
OnCalendar=Mon..Fri 13:35 UTC
OnCalendar=Mon..Fri 14:05 UTC
OnCalendar=Mon..Fri 14:35 UTC
OnCalendar=Mon..Fri 15:05 UTC
OnCalendar=Mon..Fri 15:35 UTC
OnCalendar=Mon..Fri 16:05 UTC
OnCalendar=Mon..Fri 16:35 UTC
OnCalendar=Mon..Fri 17:05 UTC
OnCalendar=Mon..Fri 17:35 UTC
OnCalendar=Mon..Fri 18:05 UTC
OnCalendar=Mon..Fri 18:35 UTC
OnCalendar=Mon..Fri 19:05 UTC
OnCalendar=Mon..Fri 19:50 UTC
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
